# systemd/paper-digest.service
[Unit]
Description=Paper Digest daily run (fetch -> digest -> trends -> publish)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/renjiu/Research/code/paperDigest
ExecStart=/home/renjiu/Research/code/paperDigest/run.sh
# A full backfill can take a couple of hours; daily incremental runs are short.
TimeoutStartSec=10800
