# systemd/paper-digest-night.timer
[Unit]
Description=Paper Digest overnight run (23:00) — digesting happens within the 23:00-08:00 window

[Timer]
OnCalendar=*-*-* 23:00:00
Persistent=true
Unit=paper-digest.service

[Install]
WantedBy=timers.target
