# systemd/paper-digest.timer
[Unit]
Description=Run Paper Digest daily at 09:00 (after the 23:00-08:00 digest window; catches up missed runs)

[Timer]
OnCalendar=*-*-* 09:00:00
Persistent=true
# Spread the trigger a little so an exact-9:00 wake isn't required.
AccuracySec=1min

[Install]
WantedBy=timers.target
