# systemd/paper-digest-verify.timer
[Unit]
Description=Verify the Paper Digest overnight run tomorrow morning (one-shot)

[Timer]
OnCalendar=2026-06-17 08:30:00
Persistent=true
Unit=paper-digest-verify.service

[Install]
WantedBy=timers.target
