# systemd/paper-digest-2026.timer
[Unit]
Description=Run the 2026 monthly compaction once, tonight at 23:00

[Timer]
OnCalendar=2026-06-16 23:00:00
Persistent=true
Unit=paper-digest-2026.service

[Install]
WantedBy=timers.target
