# systemd/paper-digest-night.timer
[Unit]
Description=Paper Digest overnight run (21:00) — digesting happens within the 21:00-06:00 window

[Timer]
OnCalendar=*-*-* 21:00:00
Persistent=true
Unit=paper-digest.service

[Install]
WantedBy=timers.target
